# config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/var/run/supervisord.pid
user=wikijs
loglevel=debug

[program:wikijs]
command=node server
autorestart=true
user=wikijs
directory=/home/wikijs
priority=100
